# config/systemd/user/usb-keyboard-detach.service
[Service]
Type=oneshot
ExecStart=/usr/bin/echo 'device has been detached'
