# config/systemd/user/usb-keyboard-attach.service
[Service]
Type=oneshot
ExecStart=/home/afk/bin/enable-keyboard
