# config/systemd/user/xscreensaver.service
[Unit]
Description=XScreensaver

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash
Restart=always

[Install]
WantedBy=default.target
