# config/systemd/user/xscreensaver.service
[Unit]
Description=XScreensaver

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash

[Install]
WantedBy=default.target
